# hub/cua-xfce/sandbox-api.conf
[program:sandbox-api]
command=/usr/local/bin/sandbox-api
user=root
autorestart=true
startsecs=0
stdout_logfile=/var/log/supervisor/sandbox-api.log
stderr_logfile=/var/log/supervisor/sandbox-api.error.log
priority=5
